# com.local.qwen3-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.qwen3-tts</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mini/Qwen3-TTS-MLX-Mac/.venv/bin/python</string>
        <string>/Users/mini/Qwen3-TTS-MLX-Mac/fastapi_server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mini/Qwen3-TTS-MLX-Mac</string>

    <!-- 开机自启 + 崩溃自动拉起 -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- 崩溃后 5 秒重启，避免死循环 -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- 日志 -->
    <key>StandardOutPath</key>
    <string>/Users/mini/Qwen3-TTS-MLX-Mac/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mini/Qwen3-TTS-MLX-Mac/logs/stderr.log</string>

    <!-- 环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
